# debian/var/lib/dpkg/info/passwd.conffiles
-rw-r--r-- 1 root root 134 Nov 18 2015 var/lib/dpkg/info/passwd.conffiles

/etc/cron.daily/passwd
/etc/default/useradd
/etc/pam.d/chfn
/etc/pam.d/chpasswd
/etc/pam.d/chsh
/etc/pam.d/newusers
/etc/pam.d/passwd
