# debian/var/lib/dpkg/info/udev.conffiles
-rw-r--r-- 1 root root 237 Mar 12 12:46 var/lib/dpkg/info/udev.conffiles

/etc/init.d/udev
/etc/init.d/udev-finish
/etc/init/udev-fallback-graphics.conf
/etc/init/udev-finish.conf
/etc/init/udev.conf
/etc/init/udevmonitor.conf
/etc/init/udevtrigger.conf
/etc/modprobe.d/fbdev-blacklist.conf
/etc/udev/udev.conf
