# debian/lib/systemd/system/networking.service.d/network-pre.conf
-rw-r--r-- 1 root root 32 Mar 12 12:46 lib/systemd/system/networking.service.d/network-pre.conf

[Unit]
After=network-pre.target
