# debian/lib/udev/rules.d/95-udev-late.rules
-rw-r--r-- 1 root root 155 Mar 12 12:46 lib/udev/rules.d/95-udev-late.rules

# do not edit this file, it will be overwritten on update

# run a command on remove events
ACTION=="remove", ENV{REMOVE_CMD}!="", RUN+="$env{REMOVE_CMD}"
